# watchdog/watchdog-proxy.timer
[Unit]
Description=Tailscale Proxy Watchdog Timer
Requires=watchdog-proxy.service

[Timer]
# Запуск через 2 хвилини після завантаження системи
OnBootSec=2min

# Запуск кожні 20 хвилин після останнього виконання
OnUnitActiveSec=20min

# Точність виконання (дозволяє системі оптимізувати використання ресурсів)
AccuracySec=1min

# Якщо система була вимкнена під час запланованого запуску,
# виконати при наступному завантаженні
Persistent=true

[Install]
WantedBy=timers.target